# common/systemd/automine_needs_reboot.path
[Unit]
Description=automine:trigger: a request from an unhealthy machine for a reboot

[Path]
PathChanged={{$AUTOMINE_ALERT_DIR}}/failed_gpus.txt

[Install]
WantedBy=paths.target
