# common/systemd/automine_reboot.path
[Unit]
Description=automine:trigger: a reboot caused by an unhealthy gpu

[Path]
PathChanged={{$HOME}}/var/automine/triggers/reboot_requests.txt

[Install]
WantedBy=paths.target
