# common/systemd/automine_check_network.timer
[Unit]
Description=automine:timer: check that the network is up after a reboot

[Timer]
OnBootSec=150

[Install]
WantedBy=timers.target
