# common/systemd/automine_reboot.service
[Unit]
Description=automine:a reboot caused by an unhealthy gpu

[Service]
ExecStart=/bin/systemctl reboot --force --force
