# common/systemd/automine_overclock.path
[Unit]
Description=automine:trigger: update the gpu overclocks

[Path]
PathChanged=/tmp/automine/automine_config.json

[Install]
WantedBy=paths.target
