# common/systemd/automine_reboot.path
[Unit]
Description=automine:trigger: a reboot caused by an unhealthy gpu

[Path]
PathChanged={{$AUTOMINE_ALERT_DIR}}/reboot_requests.txt

[Install]
WantedBy=paths.target
