# common/systemd/automine_wait_then_reboot.timer
[Unit]
Description=automine:schedule: reboot after a short delay

[Timer]
OnActiveSec=30

[Install]
WantedBy=timers.target
