# common/systemd/automine_needs_reboot.path
[Unit]
Description=automine:trigger: a request from an unhealthy gpu for a reboot

[Path]
PathChanged=%h/var/automine/triggers/failed_gpus.txt

[Install]
WantedBy=paths.target
