# common/systemd/automine_gpu_health.timer
[Unit]
Description=automine:timer: check the condition of the gpus

[Timer]
OnCalendar=*-*-* *:00,15,30,45:00

[Install]
WantedBy=timers.target
