# common/systemd/automine_reboot.service
[Unit]
Description=automine:a reboot caused by an unhealthy gpu

[Service]
ExecStart=systemctl reboot --force
