# common/systemd/automine_wait_then_overclock.timer
[Unit]
Description=automine:schedule: update overclocks after a short delay

[Timer]
OnActiveSec=15

[Install]
WantedBy=timers.target
