# common/automine_triggers.service
[Unit]
Description=Restart the screen session in response to a trigger

[Service]
ExecStart=/bin/systemctl --user restart automine
