# common/automine_reboot.service
[Unit]
Description=a reboot caused by an unhealthy gpu

[Service]
ExecStart=/sbin/reboot -f
