# common/automine_overclock.path
[Unit]
Description=automine:trigger: update the gpu overclocks

[Path]
PathChanged=/tmp/automine/overclock.json

[Install]
WantedBy=paths.target
