# common/automine_overclock.path
[Unit]
Description=trigger: update the gpu overclocks

[Path]
PathChanged=/tmp/automine/overclock.json

[Install]
WantedBy=default.target
