# common/automine_triggers.service
[Unit]
Description=automine:restart the miner because a trigger file was updated

[Service]
ExecStart=/bin/systemctl --user restart automine
