# common/automine_wait_then_reboot.timer
[Unit]
Description=scheduled: reboot after a short delay

[Timer]
OnActiveSec=30

[Install]
WantedBy=default.target
